# roles/systemd/templates/dl-git.service.j2
[Unit]
Description=Download latest git sources for the startpage & website

[Service]
ExecStart=sh "/var/cache/garuda-{{ inventory_hostname | lower }}/dl-git.sh"
