# roles/systemd/templates/goaccess.timer.j2
[Unit]
Description=Update goaccess parsed report.html

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=timers.target