# roles/systemd/files/cf-sync-ip.timer
[Unit]
Description=Get a recent list of Cloudflare IPs 

[Timer]
OnCalendar=daily
Persistent=false

[Install]
WantedBy=timers.target